# deployment/supervisor/biscuit.ini
[program:oven]
command=/home/biscuit/venv/bin/gunicorn -b localhost:8000 -w 4 oven:app
directory=/home/biscuit
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true